# services/system/update-checker/update-checker.timer
[Unit]
Description=Run HA Update Checker on weekdays at 09:00

[Timer]
OnCalendar=Mon..Fri 09:00
Persistent=false

[Install]
WantedBy=timers.target
